# systemd/orivo-sync.service
[Unit]
Description=Orivo database sync to Turso

[Service]
Type=oneshot
ExecStart=/usr/local/bin/orivo sync
