# systemd/orivo-sync.timer
[Unit]
Description=Sync Orivo todos to Turso cloud every hour

[Timer]
OnBootSec=2min
OnUnitActiveSec=1h
Persistent=true

[Install]
WantedBy=timers.target
